# webarchiv.cz/roles/manet/templates/supervisord_manet.ini.j2
[program:manet]
command=manet
autostart=true
user={{ seeder_user }}
redirect_stderr=true
